# s6-overlay/s6-rc.d/catalog-discovery/run
#!/command/execlineb -P
fdmove -c 2 1
cd /app
exec python3 -m discovery.loops.catalog_loop
